# contrib/systemd/running-process-broker-v1.service
[Unit]
Description=running-process broker v1
Documentation=https://github.com/zackees/running-process
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/running-process-broker-v1
Restart=on-failure
RestartSec=2s
KillMode=process

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
LockPersonality=yes
MemoryDenyWriteExecute=yes

RuntimeDirectory=running-process
StateDirectory=running-process
CacheDirectory=running-process
ConfigurationDirectory=running-process

ReadWritePaths=%t/running-process
ReadWritePaths=%h/.cache/running-process
ReadWritePaths=%h/.config/running-process
ReadWritePaths=%h/.local/share/running-process

[Install]
WantedBy=default.target
